# ports/sysutils/kleened/files/kleened.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: kleened
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add this line to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# kleened_enable (bool):	Set to NO by default.
#				Set it to YES to enable doormand.

. /etc/rc.subr
name=kleened
rcvar=kleened_enable
extra_commands="init dryinit"

LANG=en_US.UTF-8
LC_ALL=en_US.UTF-8

command="%%PREFIX%%/lib/kleened/bin/kleened"

pidfile="/var/run/kleened.pid"
procname="$(/usr/bin/find %%PREFIX%%/lib/kleened/erts-* -depth 1 -type d)/beam.smp"
start_cmd="${command} daemon"
init_cmd="kleened_init"
dryinit_cmd="kleened_dryinit"

load_rc_config $name
: ${kleened_enable:=no}

kleened_init()
{
    /usr/local/lib/kleened/bin/kleened eval "Kleened.Core.Config.initialize_host(%{dry_run: false})"
}

kleened_dryinit()
{
    /usr/local/lib/kleened/bin/kleened eval "Kleened.Core.Config.initialize_host(%{dry_run: true})"
}

run_rc_command "$1"
