# ports/sysutils/kleened/files/kleened.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: kleened
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add this line to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# kleened_enable (bool):	Set to NO by default.
#				Set it to YES to enable doormand.

. /etc/rc.subr

name=kleened
rcvar=kleened_enable 

start_cmd="/usr/bin/env LANG=en_US.UTF-8 LC_ALL=en_US.UTF-8 /usr/local/lib/kleened/bin/kleened daemon"
stop_cmd="/usr/bin/env LANG=en_US.UTF-8 LC_ALL=en_US.UTF-8 /usr/local/lib/kleened/bin/kleened stop"
stop_status=":"

load_rc_config $name 
: ${kleened_enable:=no} 

run_rc_command "$1"
